# vm/omv.service
[Unit]
Description=OMV Pipeline dashboard (Flask + SocketIO, port 5000)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu

# >>> VERIFIEER deze twee regels: ze moeten matchen met hoe je 'python app.py'
# >>> nu handmatig start. WorkingDirectory = de map waar app.py staat en van
# >>> waaruit je het start; de ExecStart-python = jouw venv (of 'python3').
WorkingDirectory=/home/ubuntu/omv_pipeline/v1
ExecStart=/home/ubuntu/omv_pipeline/.venv/bin/python app.py

# Portal-SSO: vertrouw de X-authentik-headers van AppPortal's nginx en sla de
# eigen login/TOTP over. Werkt samen met sso_auth.py + init_sso(app) in app.py.
Environment=AUTH_MODE=sso

# Auto-herstart bij een crash (lost het "ligt plat na een reboot"-probleem op).
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
